# contrib/systemd/ccbridge.service
[Unit]
Description=ccbridge — Claude Code hook aggregator
Documentation=https://github.com/mbocevski/ccbridge
After=graphical-session.target
PartOf=graphical-session.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=notify
ExecStart=/usr/bin/ccbridged
Restart=on-failure
RestartSec=5

WatchdogSec=60s

RuntimeDirectory=ccbridge
RuntimeDirectoryMode=0700

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ccbridge

[Install]
WantedBy=graphical-session.target
